# etc/minecraft-stop.service
[Unit]
Description=Stop Minecraft Service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl stop minecraft.service